# xorg.service
[Unit]
Description=Xorg server

[Service]
Type=notify
ExecStart=/usr/bin/xorg-launch-helper :0 -nolisten tcp -noreset -nocursor -verbose
Restart=always
RestartSec=10
