# debian/postinst
#!/bin/sh
set -e

if [ $1 = configure ]; then
	systemctl enable xorg@0.socket >/dev/null
fi

#DEBHELPER#
