# debian/prerm
#!/bin/sh

systemctl disable xorg@0.socket >/dev/null || true

#DEBHELPER#
